# trive-ai-backend.service
[Unit]
Description=TRIVE AI Backend
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/ai_copywright_trive
ExecStart=/root/ai_copywright_trive/ai_backend/venv/bin/uvicorn ai_backend.main:app --host 0.0.0.0 --port 8001
Restart=always
RestartSec=10
Environment=NODE_ENV=production
Environment=PYTHONPATH=/root/ai_copywright_trive

[Install]
WantedBy=multi-user.target 